# dist/upstart/prod.conf
# should go in /etc/init/flightfox.conf
description "Flightfox Node App"
author "flightfox"

start on startup
stop on shutdown
exec sudo NODE_ENV=production node /var/www/flightfox/live/bin/www  >> /var/www/flightfox/production.log 2>&1
