# dist/upstart/staging.conf
# should go in /etc/init/flightfox.conf
description "Flightfox Node App"
author "flightfox"

start on startup
stop on shutdown
exec sudo NODE_ENV=staging DEBUG=flightfox node /var/www/flightfox/live/bin/www  >> /var/www/flightfox/staging.log 2>&1
